# scripts/fcc-uls-downloader.service
[Unit]
Description=FCC ULS Land Mobile Database Downloader
After=network.target

[Service]
Type=oneshot
User=daryl
WorkingDirectory=/home/daryl/Python_Projects/TEA-FCC
ExecStart=/home/daryl/Python_Projects/TEA-FCC/scripts/run_download.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
